# celery.ini
[program:celery-crawler]
command=celery -A tasks worker --beat
process_name=celery-crawler
numprocs=1
autostart=true
autorestart=true
directory=/home/agendati/rabbitmq/
user=nobody