# celery.ini
[program:celery-crawler]
command=/home/agendati/rabbitmq/celery -A tasks worker --beat
process_name=celery-crawler
numprocs=1
autostart=true
autorestart=true