# mingus_server.ini
[program:mingus-server]
process_name = mingus-server-%(process_num)s
Command=python /home/agendati/backend/mingus/app.py --port=%(process_num)s
autorestart=true
autostart=true
numprocs=4
numprocs_start = 8888