# mingus_server.ini
[program:mingus-server]
command=python /home/agendati/backend/mingus/app.py
process_name=mingus_server
numprocs=1